# src/examples/port_guard_21450.conf
[program:port_guard_21450]
directory=/root/Qwen-Image/src/examples
command=/usr/bin/python3 /root/Qwen-Image/src/examples/port_guard_21450.py
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
user=root
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/port_guard_21450.out.log
stderr_logfile=/var/log/port_guard_21450.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
